# docs/ops/get2b-expire-listings.timer
[Unit]
Description=Run GET2B expire-listings every hour

[Timer]
OnCalendar=hourly
Persistent=true
Unit=get2b-expire-listings.service

[Install]
WantedBy=timers.target
